# bot.conf
[program:bot_process]
command=/root/bot/csb_retail_store/venv/bin/python -u bot.py
directory=/root/bot/csb_retail_store
stdout_logfile=/root/bot/bot_process_output.txt
redirect_stderr=true
