# bot.conf
[program:bot_process]
command=/root/bot/telebot.postgres/venv/bin/python -u bot.py
directory=/root/bot/telebot.postgres
stdout_logfile=/root/bot/bot_process_output.txt
redirect_stderr=true
